# config/flask_project.conf
[program:flask_project]
command = /home/www/env/bin/gunicorn app:app -b localhost:8000
directory = /home/www/flask_project
user = newuser