# config/flask_project.conf
[program:flask_project]
command = gunicorn app:app -b localhost:8000
directory = /home/www/flask_project